# restartor.conf
description "Restartor VNIST WEBCHECKER"
start on runlevel [2345]
stop on runlevel [06]
respawn

chdir /var/www/core/webChecker/
exec python restart.py
