# webChecker.conf
description "WebChecker VNIST version 1.0"
start on runlevel [2345]
stop on runlevel [06]
respawn

chdir /var/www/core/webChecker/
exec python main.py
